# .docker/app/rootfs/etc/services.d/workspace/run
#!/bin/bash
set -e -u -x

cd /srv/workspace
yarn start
